# data/systemd/wf-session@.target.in
[Unit]
Description=Wayfire Session Target For Instance %i
StopWhenUnneeded=true