# data/systemd/wf-session.service.in
[Unit]
Description=Wayfire Session Manager
BindsTo=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=@bindir@/wf-session-systemd
Restart=on-abnormal
RestartSec=1s
ImportCredential=wayland.*

[Install]
WantedBy=graphical-session.target
